# sources/supervisord.conf
[supervisord]

[supervisorctl]

[program:app-uwsgi]
command = /usr/bin/uwsgi --ini /data/k8sdesktop.ini

[program:nginx-app]
command = /usr/sbin/nginx
